# packaging/systemd/speculum.service
# Speculum fullscreen kiosk service — for Raspberry Pi OS Lite / headless setups
# without a desktop session. Uses `cage`, a minimal Wayland kiosk compositor, to
# run the mirror fullscreen on the console.
#
# Install:
#   sudo apt install cage          # (Debian/Ubuntu) or: sudo pacman -S cage
#   # edit User= below to the account that owns the display, then:
#   sudo cp speculum.service /etc/systemd/system/    # (deb/rpm ship it under
#   #   /opt/speculum/lib/app/resources/; Arch installs it automatically)
#   sudo systemctl enable --now speculum.service
#
# On Arch the unit is already at /usr/lib/systemd/system/speculum.service.

[Unit]
Description=Speculum smart-mirror (fullscreen kiosk)
After=systemd-user-sessions.service network-online.target
Wants=network-online.target

[Service]
# The account that owns seat0 / the display. Change to your user (e.g. "pi").
User=pi
# Open a PAM session so pam_systemd creates /run/user/<uid> and exports
# XDG_RUNTIME_DIR — cage needs it for the Wayland socket. Without this a system
# unit has no user session and cage aborts:
#   [ERROR] [../cage.c:280] XDG_RUNTIME_DIR is not set in the environment
# Uses our own minimal stack (/etc/pam.d/speculum) rather than "login" so we
# don't drag in pam_lastlog/pam_motd — pam_lastlog is gone on newer distros and
# logs "adding faulty module". deb/rpm/Arch install the file; for the manual
# copy above, also: sudo cp packaging/pam/speculum /etc/pam.d/speculum
PAMName=speculum
# Compose Desktop renders through AWT, so it needs Xwayland (DISPLAY=:1). Don't
# force the *compositor* to software (WLR_RENDERER=pixman) — that strips GL from
# wlroots and makes Xwayland's startup fail:
#   [xwayland/server.c] Xwayland startup failed
#   java.awt.AWTError: Can't connect to X11 window server using ':1'
# Instead leave cage/Xwayland on the GPU and force only Skiko (the app's
# renderer) to software raster, which avoids GL crashes on old Intel GPUs.
# JAVA_TOOL_OPTIONS is picked up automatically by the bundled JRE.
# Drop this on hardware with working GPU acceleration for better performance.
Environment=JAVA_TOOL_OPTIONS=-Dskiko.renderApi=SOFTWARE
# cage launches a minimal Wayland session and runs Speculum fullscreen in it.
ExecStart=/usr/bin/cage -s -- /opt/speculum/bin/speculum
Restart=on-failure
RestartSec=3
# Bind to the console so cage can take the seat.
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target